# inu-logistics-backend.service
[Unit]
Description=INU Logistics Backend Server
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/inu/INU_final
Environment=PATH=/home/inu/INU_final/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/inu/INU_final/venv/bin/python -m backend.app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give the service time to start up
TimeoutStartSec=60

# Environment variables
Environment=PYTHONPATH=/home/inu/INU_final
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target 